# etc/xdg/autostart/lxqt-copyq-autostart.desktop
[Desktop Entry]
Type=Application
Name=CopyQ
TryExec=copyq
Exec=copyq
OnlyShowIn=LXQt;
X-LXQt-Need-Tray=true

#TRANSLATIONS_DIR=translations
